# root-arm/etc/s6-overlay/s6-rc.d/init-calibre-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cat <<-EOF
    ********************************************************
    ********************************************************
    *                                                      *
    *                         !!!!                         *
    *         This Calibre image does not support          *
    *        ARM due to a lack of upstream packages        *
    *                                                      *
    *                                                      *
    ********************************************************
    ********************************************************
EOF
